# greenpulse/rootfs/etc/services.d/greenpulse/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting GreenPulse Application..."
# We need to wait for MariaDB to be ready.
# A simple sleep loop or wait-for-it approach.
bashio::log.info "Waiting for database..."
sleep 30

exec python3 /app/src/main.py
